# deploy/pocketbase.service
[Unit]
Description=PocketBase (LinktreeSRP)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/LinktreeSRP/pb
ExecStart=/root/LinktreeSRP/pb/pocketbase serve --http=127.0.0.1:8090
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
